# deploy/kiosk.service
# deploy/kiosk.service
# Launches the Tauri dashboard app after the Express bridge is ready.
#
# Install: sudo cp deploy/kiosk.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable kiosk
#          sudo systemctl start kiosk

[Unit]
Description=Maverick Telemetry — Kiosk (Tauri)
Documentation=https://github.com/AlexTs-dev/maverick-telemetry-hub
After=express_bridge.service graphical.target
Requires=express_bridge.service

[Service]
Type=simple
User=pi

# Wait for Express bridge to be ready before opening the window.
# Retries curl every second for up to 30 seconds.
ExecStartPre=/bin/bash -c '\
  for i in $(seq 1 30); do \
    curl -sf http://localhost:3000/api/health && exit 0; \
    sleep 1; \
  done; \
  exit 1'

ExecStart=/home/pi/maverick-telemetry-hub/deploy/kiosk-start.sh

Restart=on-failure
RestartSec=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=kiosk

[Install]
WantedBy=graphical.target
